# systemd/ambiance-mpd.service
[Unit]
Description=ambiance-amplipi :: mpd (radio player -> ch0)
Documentation=file:/home/openhab/docs/audio/APPLIANCE_RADIO_PLAN.md
After=network.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /home/pi/ambiance-amplipi/run/playlists
ExecStart=/usr/bin/mpd --no-daemon /home/pi/ambiance-amplipi/config/mpd.conf
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
